# doc/ansible/templates/sysmon-client.service.j2
[Unit]
Description=Sysmon Client Service

[Service]
TimeoutStartSec=3
Restart=always
ExecStart=/opt/sysmon/client/bin/client -s {{server_url}}

[Install]
WantedBy=default.target
